# daemon/tinyclaw-discord.service
[Unit]
Description=TinyClaw Trading — Discord Client
After=network.target tinyclaw-trading.service
Requires=tinyclaw-trading.service

[Service]
Type=simple
User=strategist
Group=strategist
WorkingDirectory=/opt/ArkNode-AI/projects/trading-signal-ai/tinyclaw-trading
EnvironmentFile=/opt/ArkNode-AI/projects/trading-signal-ai/tinyclaw-trading/.env
Environment=HOME=/home/strategist
Environment=PATH=/home/strategist/.bun/bin:/home/strategist/.local/bin:/usr/local/bin:/usr/bin:/bin

ExecStart=/home/strategist/.bun/bin/bun run src/channels/discord-client.ts

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tinyclaw-discord

NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
